# dd.conf
[program:dd]
command=python3 /home/curator2/dd/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/dd/test.err.log
stdout_logfile=/home/curator2/dd/test.out.log